# devopsproject/root/etc/supervisor/conf.d/devopsproject.ini
[program:devopsproject]
command=python /opt/bin/devopsproject.py
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/devopsproject.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/devopsproject-error.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
